# devops/systemd/deploy-webhook.service
[Unit]
Description=MegiLance GitHub Webhook Deployer
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/deploy/webhook
Environment=NODE_ENV=production
# WEBHOOK secret should be set via a secure env file (see instructions)
EnvironmentFile=/etc/megilance_deploy.env
ExecStart=/usr/bin/node /opt/deploy/webhook/index.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
